# infra/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24
autorestart=true
priority=100

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=200

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport %(ENV_VNC_PORT)s -nopw
autorestart=true
priority=300

[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autorestart=true
priority=400

[program:xrce-agent]
command=/usr/local/bin/micro-xrce-dds-agent udp4 -p 8888
autorestart=true
priority=500
user=simuser

[program:px4]
command=/bin/bash -c "cd /home/simuser/PX4-Autopilot && make px4_sitl %(ENV_PX4_MODEL)s"
environment=DISPLAY=":99",HEADLESS="%(ENV_HEADLESS)s",PX4_MODEL="%(ENV_PX4_MODEL)s"
autorestart=false
priority=600
user=simuser
